# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:blog]
directory=/code
command=jekyll serve --host 0.0.0.0
